# scripts/auto_ingest_sermon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pastordavepro.auto-ingest-sermon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/glennheistand/Projects/pastordave-pro/scripts/auto_ingest_sermon.py</string>
        <string>--limit</string>
        <string>1</string>
        <string>--church-id</string>
        <string>new-horizon-champaign</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/glennheistand/Projects/pastordave-pro</string>
    
    <!-- Note: OPENAI_API_KEY and ANTHROPIC_API_KEY must be set in ~/.bash_profile or ~/.zshrc -->
    <!-- The script reads them from the environment when launchd runs -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <!-- Run every Sunday at 8:00 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/auto_ingest_sermon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/auto_ingest_sermon_error.log</string>
</dict>
</plist>
